# tools/pkg-scripts/deb/debian-0.9/postinst
#!/bin/sh

python3 -m pip install pyparsing pylibconfig2

systemctl enable sx-network@default.service
systemctl enable sx-core@default.service